# sv/strider/run
#!/bin/bash

sv start mongodb || exit 1

while ! nc -vz localhost 27017;do sleep 3; done

source /etc/envvars
env

exec 2>&1
exec /strider/bin/strider
